# config/systemd/piclaw-telegram.service
[Unit]
Description=PiClaw Telegram bridge (forwards to OpenClaw agent)
After=network-online.target openclaw.service
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/openclaw/src

# Token: create bot via @BotFather, then set here or in /opt/openclaw/config/telegram.env
EnvironmentFile=-/opt/openclaw/config/telegram.env
Environment=OPENCLAW_CHAT_URL=http://127.0.0.1:3100/chat

ExecStart=/usr/bin/node /opt/openclaw/src/telegram-bridge.js
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=piclaw-telegram

[Install]
WantedBy=multi-user.target
